# services/sensor.service
[Unit]
Description=Sensor Data Uploader
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/pi/src/sensor_uploader.py
WorkingDirectory=/home/pi/src/
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
